# datasets/scenario4_supply_chain/metadata.json
{
  "scenario": {
    "id": 4,
    "codename": "Supply Chain Phantom",
    "name": "Operation Supply Chain Phantom",
    "description": "Multi-vector supply chain attack via dependency confusion, DNS tunneling, data exfiltration, and anti-forensics",
    "difficulty": "advanced_only",
    "version": "1.0"
  },
  "timeframe": {
    "start": "2026-03-05T00:00:00Z",
    "end": "2026-03-06T23:59:59Z",
    "attack_start": "2026-03-05T02:14:33Z",
    "attack_end": "2026-03-06T06:00:00Z",
    "timezone": "UTC"
  },
  "infrastructure": {
    "victims": [
      {
        "hostname": "web-srv",
        "ip": "172.26.0.30",
        "wazuh_agent_id": "001"
      },
      {
        "hostname": "dc-srv",
        "ip": "172.26.0.32",
        "wazuh_agent_id": "002"
      },
      {
        "hostname": "lnx-srv",
        "ip": "172.26.0.31",
        "wazuh_agent_id": "003"
      }
    ],
    "attacker": {
      "c2_domains": [
        "cdn-analytics.cloud-metrics.net",
        "cdn-static.cloud-metrics.net"
      ],
      "dns_tunnel_domain": "cdn-analytics.cloud-metrics.net",
      "exfil_endpoint": "https://cdn-static.cloud-metrics.net/api/v2/upload"
    },
    "malicious_package": {
      "name": "wazuhbots-utils",
      "registry": "pip",
      "backdoor_path": "/usr/local/lib/python3.10/dist-packages/wazuhbots_utils/.config/svc_update.py"
    }
  },
  "flags": {
    "pup_01_hosts_infected": "3",
    "pup_02_package_manager": "pip",
    "pup_03_attack_date": "2026-03-05",
    "pup_04_service_name": "svc_update.service",
    "pup_05_first_host": "web-srv",
    "pup_06_dns_count_websrv": "40",
    "pup_07_backdoor_script": "svc_update.py",
    "pup_08_logrotate_config": "/etc/logrotate.d/wazuh",
    "pup_09_package_version": "1.3.7",
    "pup_10_pip_user": "root",
    "hunter_01_infection_spread": "390",
    "hunter_02_dcsrv_infection": "2026-03-05T06:31:07Z",
    "hunter_03_dns_base_domain": "cdn-analytics.cloud-metrics.net",
    "hunter_04_exfil_domain": "cdn-static.cloud-metrics.net",
    "hunter_05_websrv_exfil_mb": "312",
    "hunter_06_staging_dirs": "/tmp/.cache/|/dev/shm/",
    "hunter_07_new_logrotate": "rotate_0_maxage_1",
    "hunter_08_logrotate_force": "logrotate -f /etc/logrotate.d/wazuh",
    "hunter_09_lnxsrv_infection": "2026-03-05T08:45:22Z",
    "hunter_10_process_parent": "pip",
    "alpha_01_malicious_package": "wazuhbots-utils",
    "alpha_02_backdoor_path": "/usr/local/lib/python3.10/dist-packages/wazuhbots_utils/.config/svc_update.py",
    "alpha_03_dns_tunnel_domain": "cdn-analytics.cloud-metrics.net",
    "alpha_04_first_affected_host": "web-srv",
    "alpha_05_dns_total": "120",
    "alpha_06_supply_chain_mitre": "T1195.001",
    "alpha_07_dcsrv_exfil_mb": "287",
    "alpha_08_backdoor_uid": "0",
    "alpha_09_chunk_pattern": "chunk_NNNN.enc",
    "alpha_10_original_logrotate": "rotate_7_maxage_30",
    "alpha_11_systemd_unit": "/etc/systemd/system/svc_update.service",
    "alpha_12_package_hash": "c7a5b3d9e2f14680ab91cd3e4f567890123456789abcdef0123456789abcdef0",
    "fenrir_01_multi_host_timeline": "web-srv:2026-03-05T02:14:33Z|dc-srv:2026-03-05T06:31:07Z|lnx-srv:2026-03-05T08:45:22Z",
    "fenrir_02_exfil_endpoint_volume": "https://cdn-static.cloud-metrics.net/api/v2/upload:847",
    "fenrir_03_logrotate_manipulation": "/etc/logrotate.d/wazuh:rotate_0_maxage_1:web-srv,lnx-srv",
    "fenrir_04_exfil_breakdown": "web-srv:312|dc-srv:287|lnx-srv:248",
    "fenrir_05_infection_delays": "256|134",
    "fenrir_06_anti_forensics_hosts": "lnx-srv|web-srv",
    "fenrir_07_subdomain_length": "32",
    "fenrir_08_c2_domains": "cdn-analytics.cloud-metrics.net|cdn-static.cloud-metrics.net",
    "fenrir_09_mitre_chain": "T1041|T1059.006|T1070|T1071.004|T1074.001|T1195.001|T1543.002",
    "fenrir_10_backdoor_md5": "e99a18c428cb38d5f260853678922e03"
  },
  "mitre_techniques": [
    "T1195.001",
    "T1059.006",
    "T1071.004",
    "T1543.002",
    "T1074.001",
    "T1041",
    "T1070"
  ],
  "total_alerts": 600,
  "total_archives": 9200
}