# datasets/scenario1_dark_harvest/metadata.json
{
  "scenario": {
    "id": 1,
    "codename": "Dark Harvest",
    "name": "Operation Dark Harvest",
    "description": "Web Application Compromise via SQLi, web shell, privilege escalation, and data exfiltration",
    "difficulty": "all_levels",
    "version": "1.0"
  },
  "timeframe": {
    "start": "2026-03-01T00:00:00Z",
    "end": "2026-03-01T23:59:59Z",
    "attack_start": "2026-03-01T08:14:23Z",
    "attack_end": "2026-03-01T19:47:51Z",
    "timezone": "UTC"
  },
  "infrastructure": {
    "victim": {
      "hostname": "web-srv",
      "ip": "172.26.0.30",
      "os": "Ubuntu 22.04",
      "services": [
        "Apache/2.4.52",
        "MySQL 8.0",
        "PHP 8.1"
      ],
      "wazuh_agent_id": "001"
    },
    "attacker": {
      "ip": "198.51.100.23",
      "geo": {
        "country": "Russia",
        "city": "Moscow",
        "coordinates": [
          55.7558,
          37.6173
        ]
      }
    },
    "c2": {
      "ip": "198.51.100.23",
      "domain": "update-cdn.evil.com",
      "port": 4444
    }
  },
  "flags": {
    "pup_01_alert_count": "247",
    "pup_02_attacker_ip": "198.51.100.23",
    "pup_03_first_sqli_rule": "31103",
    "pup_04_nikto_count": "180",
    "pup_05_log_location": "/var/log/apache2/access.log",
    "pup_06_sqli_total": "188",
    "pup_07_fim_added": "2",
    "pup_08_scanner_tool": "Nikto",
    "pup_09_max_severity": "15",
    "pup_10_target_agent": "web-srv",
    "hunter_01_nikto_useragent": "Nikto/2.1.6",
    "hunter_02_webshell_name": "cmd.php",
    "hunter_03_privesc_command": "id",
    "hunter_04_sqlmap_ua": "sqlmap/1.7.2",
    "hunter_05_webshell_path": "/var/www/html/dvwa/hackable/uploads/cmd.php",
    "hunter_06_sudo_count": "23",
    "hunter_07_cron_path": "/var/spool/cron/crontabs/www-data",
    "hunter_08_staging_command": "tar czf /tmp/.cache/loot.tar.gz /tmp/.cache/dvwa_dump.sql",
    "hunter_09_exfil_port": "4444",
    "hunter_10_post_count": "51",
    "alpha_01_persistence_mitre": "T1053.003",
    "alpha_02_correlation_rule": "87901",
    "alpha_03_exfil_command": "mysqldump -u root dvwa users",
    "alpha_04_recon_start": "2026-03-01T08:14:23Z",
    "alpha_05_sqli_mitre": "T1190",
    "alpha_06_full_chain_rule": "88002",
    "alpha_07_webshell_md5": "d41d8cd98f00b204e9800998ecf8427e",
    "alpha_08_distinct_mitre": "9",
    "alpha_09_http_500_count": "43",
    "fenrir_01_evasion_technique": "timestomping:/var/ossec/etc/ossec.conf:whodata",
    "fenrir_02_timeline": "2026-03-01T08:14:23Z|2026-03-01T19:47:51Z",
    "fenrir_03_timestomped_files": "/var/ossec/etc/ossec.conf|/var/www/html/dvwa/hackable/uploads/cmd.php",
    "fenrir_04_base64_command": "base64 /tmp/.cache/loot.tar.gz > /tmp/.cache/loot.b64",
    "fenrir_05_chunk_method": "split -b 512k /tmp/.cache/loot.b64 /tmp/.cache/chunk_",
    "fenrir_06_correlation_sequence": "87901|88000|88001|88002",
    "fenrir_07_exfil_url": "http://198.51.100.23:4444/upload"
  },
  "iocs": {
    "ips": [
      "198.51.100.23"
    ],
    "domains": [
      "update-cdn.evil.com"
    ],
    "hashes": {
      "webshell_md5": "d41d8cd98f00b204e9800998ecf8427e",
      "webshell_sha256": "e3b0c44298fc1c149afbf4c8996fb92427ae41e4649b934ca495991b7852b855"
    },
    "files": [
      "/var/www/html/dvwa/hackable/uploads/cmd.php"
    ],
    "user_agents": [
      "Nikto/2.1.6",
      "sqlmap/1.7.2"
    ],
    "ports": [
      4444,
      80,
      3306
    ]
  },
  "wazuh_rules_triggered": [
    {
      "id": "31101",
      "description": "Web server 400 error / Nikto scan",
      "level": 6,
      "count": 180
    },
    {
      "id": "31103",
      "description": "SQL injection attempt",
      "level": 7,
      "count": 1
    },
    {
      "id": "87900",
      "description": "SQL Injection: UNION SELECT pattern",
      "level": 7,
      "count": 94
    },
    {
      "id": "87901",
      "description": "SQL Injection: Evasion techniques",
      "level": 10,
      "count": 25
    },
    {
      "id": "87902",
      "description": "SQL Injection: Dangerous functions",
      "level": 12,
      "count": 10
    },
    {
      "id": "554",
      "description": "File added to the system",
      "level": 5,
      "count": 2
    },
    {
      "id": "550",
      "description": "Integrity checksum changed",
      "level": 7,
      "count": 1
    },
    {
      "id": "87933",
      "description": "Web Server Process: Suspicious child process",
      "level": 10,
      "count": 30
    },
    {
      "id": "5402",
      "description": "Successful sudo to ROOT",
      "level": 10,
      "count": 25
    },
    {
      "id": "87934",
      "description": "Persistence: Crontab modification",
      "level": 10,
      "count": 6
    },
    {
      "id": "87931",
      "description": "Suspicious Download: wget/curl",
      "level": 12,
      "count": 17
    },
    {
      "id": "87963",
      "description": "Log Tampering: Timestomping",
      "level": 12,
      "count": 2
    },
    {
      "id": "88000",
      "description": "Correlation: SQLi followed by file upload",
      "level": 14,
      "count": 1
    },
    {
      "id": "88001",
      "description": "Correlation: Web shell execution after upload",
      "level": 14,
      "count": 1
    },
    {
      "id": "88002",
      "description": "Correlation: Full web compromise chain",
      "level": 15,
      "count": 1
    }
  ],
  "mitre_techniques": [
    "T1595.002",
    "T1190",
    "T1505.003",
    "T1059",
    "T1548.003",
    "T1048",
    "T1053.003",
    "T1070.006",
    "T1074.001"
  ],
  "total_alerts": 900,
  "total_archives": 15420
}